# deploy/botcreator.service
[Unit]
Description=Helles-Galerie Bot Creator Web UI
After=network.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/PicoGallery/Chatbot/BotCreator
ExecStart=/usr/bin/python3 /home/admin/PicoGallery/Chatbot/BotCreator/web_creator.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
